# etc/systemd/system/nginx.service.d/local.conf
# Replace this environment variable with your own domain
# It is used by ngnix config to know what it is listening on.
[Service]
Environment="SERVER_NAME=your_domain.com"
